# layout/DEBIAN/postinst
#!/bin/sh

# Verify supersling permissions
chown root:wheel /usr/libexec/zebra/supersling
chmod 6755 /usr/libexec/zebra/supersling

# Setup relaunch Launch Daemon
launchctl load /Library/LaunchDaemons/xyz.willy.Zebra.Relaunch.plist